# etc/cron/crontab
* * * * * /gost/gost -app=script -name=test
